# arch/.config/systemd/user/hyprsunset-day.timer
[Unit]
Description=Hyprsunset Day Timer

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target